# backend/systemd/peak-music.service
[Service]
ExecStart=/usr/bin/peak-music
Environment="ADDR=0.0.0.0:80"
